# ansible/roles/lb/templates/lb.service
[Unit]
Description=Load Balancer

[Service]
TimeoutStartSec=0
ExecStart={{lb_dir}}/fabio
